# ubumirror.crontab
# https://wiki.ubuntu.com/Mirrors
# Update archive mirror every six hours.
03 */6 * * * ubumirror /usr/bin/ubuarchive >/dev/null
# Update releases mirror every four hours.
03 */4 * * * ubumirror /usr/bin/uburelease >/dev/null
# Update cdimage mirror every twenty-four hours.
03 11 * * * ubumirror /usr/bin/ubucdimage >/dev/null
